# claude-mem/com.kryss.claude-mem.auto-cleanup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.kryss.claude-mem.auto-cleanup</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/kryss/projects/macbook-tools/claude-mem/auto-cleanup.sh</string>
    </array>

    <!-- 每天午夜 0:00 运行 -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>0</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/kryss/.claude-mem/auto-cleanup.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/kryss/.claude-mem/auto-cleanup.stderr.log</string>

    <key>KeepAlive</key>
    <false/>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
