# claude-mem/com.kryss.claude-mem.monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.kryss.claude-mem.monitor</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/kryss/projects/macbook-tools/claude-mem/monitor.sh</string>
    </array>

    <key>StartInterval</key>
    <integer>86400</integer>

    <!-- 每天 9 点运行（可选，用 StartCalendarInterval 替代 StartInterval） -->
    <!--
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    -->

    <key>StandardOutPath</key>
    <string>/Users/kryss/.claude-mem/monitor.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/kryss/.claude-mem/monitor.stderr.log</string>

    <key>KeepAlive</key>
    <false/>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
